# zsh/.config/zsh/profile.d/0path.zsh
export PATH="${KREW_ROOT:-$HOME/.krew}/bin:$HOME/bin/roc:$PATH:$HOME/.zigup:$HOME/go/bin:/Applications/Docker.app/Contents/Resources/bin"
