# zsh/.config/zsh/profile.d/alias.zsh
alias ls='ls --color'
alias vim='nvim'
alias vi='nvim'
alias ll='ls -lart'
