# zsh/.config/zsh/profile.d/direnv.zsh
export LESS="-XFR"
eval "$(direnv hook zsh)"
