# zsh/.config/zsh/profile.d/options.zsh
fpath=(/usr/local/share/zsh-completions ~/.config/zsh/functions $fpath)
HISTSIZE=100000
SAVEHIST=100000
HISTFILE=~/.zsh_history
unsetopt menucomplete
unsetopt automenu
setopt interactivecomments
unsetopt histverify
setopt histsubstpattern
setopt sharehistory
setopt inc_append_history
setopt share_history
