# zsh/.config/zsh/profile.d/options.zsh
fpath=(/usr/local/share/zsh-completions ~/.config/zsh/functions $fpath)
unsetopt menucomplete
unsetopt automenu
setopt interactivecomments
unsetopt histverify
