# zsh/.config/zsh/profile.d/path.zsh
export PATH="$HOME/bin:$HOME/bin:/usr/local/bin:/bin:/usr/local/sbin:/sbin:/usr/bin:/usr/local/sbin:/usr/local/MacGPG2/bin:/home/walexander/perl5/bin"
