# zsh/.config/zsh/profile.d/powerline.zsh
#powerline_dir=/usr/local/lib/python2.7/site-packages/powerline/bindings
#powerline_init_file=$powerline_dir/zsh/powerline.zsh
#. ~/.p10k.zsh

#if [ -f $powerline_init_file ]; then
##	. $powerline_init_file
#	source $(brew --prefix)/opt/powerlevel10k/powerlevel10k.zsh-theme
#fi
