# zsh/.config/zsh/profile.d/vimrc.zsh
export MYVIMRC=$HOME/.config/vim/vimrc
export EDITOR=vim
